# legal-music-bot.service
[Unit]
Description=Legal Music Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/legal-music
ExecStart=/usr/bin/python3 %h/legal-music/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=append:%h/legal-music/output/bot.log
StandardError=append:%h/legal-music/output/bot.log
EnvironmentFile=%h/legal-music/.env

[Install]
WantedBy=default.target
